# assistant/systemd/nvidia-watchdog.timer
[Unit]
Description=MindHome NVIDIA GPU Watchdog (alle 300s)

[Timer]
OnBootSec=30s
OnUnitActiveSec=300s
AccuracySec=5s

[Install]
WantedBy=timers.target
